# software/var_lib_smartmold/sys.files/systemd/downtime.service
[Unit]
Description=Downtime Logger

[Service]
Type=simple
StandardError=journal
User=smartmold
ExecStart=/var/lib/smartmold/scripts/downtime.sh
