# software/var_lib_smartmold/sys.files/systemd/mold-sweat.service
[Unit]
Description=Mold Sweat Logger

[Service]
Type=simple
StandardError=journal
User=root
ExecStart=/var/lib/smartmold/scripts/moldsweat-data.sh
