# software/var_lib_smartmold/sys.files/systemd/mold-sweat-create-graphs.service
[Unit]
Description=Smart Mold mold sweat graph creator

[Service]
Type=simple
StandardError=journal
User=smartmold
ExecStart=/var/lib/smartmold/scripts/create-graphs.sh
